# deployment/systemd/totem-frontend.service
[Unit]
Description=Totem Frontend
After=network.target totem-backend.service
Wants=totem-backend.service

[Service]
Type=simple
User=DEPLOY_USER
WorkingDirectory=PROJECT_ROOT/apps/frontend
Environment="NODE_ENV=production"
Environment="BACKEND_URL=http://localhost:3001"
Environment="PORT=80"
Environment="HOST=0.0.0.0"
Environment="ORIGIN=http://192.168.5.25"
ExecStart=/usr/local/bin/bun build/index.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
